# supervisor/4-ttyd.conf
[program:ttyd]
command="/usr/local/bin/ttyd -p 7681 bash -x"
numprocs=1
startsecs=0
priority=11
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autostart = true
autorestart=true
